# config/systemd/user/desktop-battery-monitor.service
[Unit]
Description=Desktop Battery Monitor Check
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/share/desktop/bin/desktop-battery-monitor
Environment=DISPLAY=:0
